# deploy/files/lib/systemd/system/growatt.target
[Unit]
Description=Acquire Growatt data every 5 minutes (target)
StopWhenUnneeded=yes
